# ops/hetzner/age-worker.service.d/validate-zernio.conf
[Service]
ExecStartPre=/bin/bash /opt/age/ops/hetzner/validate-zernio-key.sh /opt/age/.env
